# build/root/etc/services.d/trackarr/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid abc /app/trackarr/trackarr --config-dir="${CONFIG_DIR}"
